# claude-notify-serve.service
[Unit]
Description=Claude Code Telegram Notification Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m hookline --serve
Environment=PYTHONPATH=__HOOKS_DIR__
Environment=TELEGRAM_BOT_TOKEN=__BOT_TOKEN__
Environment=TELEGRAM_CHAT_ID=__CHAT_ID__
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
